# backend/ai-learning.service
[Unit]
Description=Healzy AI Continuous Learning System (DISABLED)
Documentation=https://healzy.uz/docs/ai-learning
After=network.target mysql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/healzy.app/backend
Environment="PATH=/var/www/healzy.app/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/var/www/healzy.app/backend"
Environment="PYTHONUNBUFFERED=1"
Environment="ENABLE_AI=false"

# Команда запуска
ExecStart=/var/www/healzy.app/venv/bin/python ai_continuous_learning.py

# Перезапуск при сбоях отключен для предотвращения перезагрузок сервера
Restart=no
RestartSec=10

# Лимиты ресурсов
CPUQuota=50%
MemoryLimit=2G
MemorySwapMax=0

# Логирование
StandardOutput=append:/var/log/healzy/ai-learning.log
StandardError=append:/var/log/healzy/ai-learning-error.log

# Таймауты
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target 