# backend/vrachi-backend.service
[Unit]
Description=VrachiAP Backend Service
After=network.target

[Service]
User=www-data
WorkingDirectory=/var/www/healzy.app/backend
Environment="PATH=/var/www/healzy.app/backend/venv/bin"
Environment="PYTHONPATH=/var/www/healzy.app/backend"
Environment="ENABLE_AI=false"
ExecStart=/var/www/healzy.app/backend/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target